# minipc/router-setup/roles/dns/templates/bind-rpz-update.service.j2
# Managed by Ansible — do not edit manually
[Unit]
Description=Update Bind9 RPZ blocklist (porn + gambling)
After=network-online.target named.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/update-bind-rpz
