# minipc/router-setup/roles/captive_portal/templates/wan-check.service.j2
# Managed by Ansible — do not edit manually
[Unit]
Description=WAN health check — toggle offline mode (RPZ + nginx)
After=named.service nginx.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wan-check.sh
